# systemd/sensible.path
[Unit]
Description=Sensible Queue Watcher
Documentation=https://github.com/keithy/sensible

[Path]
DirectoryNotEmpty=/var/lib/sensible/tasks/pending
Unit=sensible.service

[Install]
WantedBy=multi-user.target
